# deploy/liubai.service
# systemd — 留白 API（/www/liubai）
# sudo cp deploy/liubai.service /etc/systemd/system/liubai.service
# sudo systemctl daemon-reload
# sudo systemctl enable --now liubai

[Unit]
Description=留白 Blog API (FastAPI)
Documentation=file:/www/liubai/DEPLOY.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data

WorkingDirectory=/www/liubai/backend
EnvironmentFile=-/etc/blog/env

ExecStart=/www/liubai/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000

Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# DATA_DIR=/www/liubai 时需可写 db.sqlite3 与 uploads/
ReadWritePaths=/www/liubai
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
